# etc/supervisord.conf
[supervisord]
nodaemon=true

[program:xx-net]
command = /XX-Net/start
autostart = true
autorestart = true

[program:miredo]
command = miredo -f
autostart = true
autorestart = true